# server/doc/isolated_foo.service
[Unit]
Requires=isolated_network.service
After=isolated_network.service

[Service]
NetworkNamespacePath=/run/netns/isolator

